# src/cloud-api-adaptor/podvm/files/etc/systemd/system/run-kata\x2dcontainers-image-overlay.mount
[Unit]
Description=Mount unit for /run/kata-containers/image/overlay
Before=kata-agent.service

[Mount]
What=tmpfs
Where=/run/kata-containers/image/overlay
Type=tmpfs

[Install]
WantedBy=multi-user.target
